# cookbooks/vertx/templates/default/supervisor-service.conf.erb
[program:<%= @mod_conf['mod'] %>]
command=<%= @vertx[:mods_dir] + "/" + @mod_conf['mod'] + ".run" %>
process_name=<%= @mod_conf['mod'] %>
autostart=true
autorestart=true
stopsignal=QUIT
user=<%= @vertx[:user] %>
group=<%= @vertx[:group] %>

